# crontab
# Run every minute
*/1 * * * * yarn run jobs:process-voicemail:production